# services/build/worker.cron
* * * * * root cd /shared/src/collecting_society_worker/collecting_society_worker && /opt/venv/bin/python repro.py all